# dockerfiles/queues/laravel-echo-worker.conf
[program:laravel-echo-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/laravel-echo-server start --dir=/home/www/signing/
autostart=true
autorestart=true
numprocs=1
