# dockerfiles/queues/laravel-echo-worker.conf
[program:laravel-echo-worker]
process_name=%(program_name)s_%(process_num)02d
command=/home/www/signing/laravel-echo-server start
autostart=true
autorestart=true
numprocs=2
