# dockerfiles/queues/wellsail-worker.conf
[program:wellsail-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /home/www/signing/artisan queue:work redis
autostart=true
autorestart=true
numprocs=2
